# autoinstall/mi325xr/common/firstboot.target
[Unit]
Description=First Boot Initialization Target
ConditionPathExists=!/var/lib/firstboot.done

[Install]
WantedBy=multi-user.target
